# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/aix-system artisan schedule:run >> /dev/null 2>&1
